# scraper/crontab
0 3 * * * root cd /app && python scraper.py >> /var/log/cron.log 2>&1

